# jim/root/var/lib/dpkg/info/nymea-networkmanager.list
/.
/etc
/etc/nymea
/etc/nymea/nymea-networkmanager.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/nymea-networkmanager.service
/usr
/usr/bin
/usr/bin/nymea-networkmanager
/usr/share
/usr/share/doc
/usr/share/doc/nymea-networkmanager
/usr/share/doc/nymea-networkmanager/changelog.gz
/usr/share/doc/nymea-networkmanager/copyright
